# config/autostart/sxhkd.desktop
[Desktop Entry]
Name=sxhkd
Comment=Simple X hotkey daemon
Exec=/usr/local/bin/sxhkd
Terminal=false
Type=Service
Hidden=false

